# srcs/influx/supervisord.conf
;Явно сообщаем supervisord что процессы запускаются не в фоновом режиме

[supervisord]
nodaemon=true
user = root

; Запускаем grafana
[program:grafana]
command= sh -c "grafana-server && kill -s SIGTERM $(cat supervisord.pid)"

; Запускаем influxd
[program:influxdb]
command= sh -c "influxd run -config /etc/influxdb.conf && kill -s SIGTERM $(cat supervisord.pid)"
